# configs/systemd/mcp-annet-oil.service
[Unit]
Description=MCP Annet Oil - Model Context Protocol Server for Annet Oil
Documentation=https://github.com/yourusername/annet-oil
After=network-online.target annet-oil.service
Wants=network-online.target annet-oil.service

[Service]
Type=simple
User=annet
Group=annet
WorkingDirectory=/opt/annet-oil/mcp-annet-oil
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"

# Environment variables for MCP server
Environment="ANNET_OIL_API_URL=http://localhost:8181"
Environment="ANNET_OIL_API_TOKEN=your-secret-token"
Environment="MCP_SERVER_NAME=annet-oil"

# Node.js memory settings
Environment="NODE_OPTIONS=--max-old-space-size=2048"

# Start command
# Dependencies should be installed manually before starting service
ExecStart=/usr/local/bin/node dist/index.js

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/annet-oil/mcp-annet-oil

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryLimit=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-annet-oil

[Install]
WantedBy=multi-user.target